# deploy/sdm-worker.service
# /etc/systemd/system/sdm-worker.service
#
# SDM Job Queue Worker — runs `flask sdm worker` as a system service.
#
# Install on the server:
#   sudo cp deploy/sdm-worker.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable sdm-worker
#   sudo systemctl start sdm-worker
#
# View logs:
#   sudo journalctl -u sdm-worker -f
#
# Stop (graceful — the worker finishes the current task before exiting):
#   sudo systemctl stop sdm-worker

[Unit]
Description=SDM Occupancy Model Worker (Flask job queue)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
# Set these to the user/group and path used by your deployment.
User=www-data
Group=www-data
WorkingDirectory=/var/www/biomon

# Activate the venv and start the worker via the flask CLI
ExecStart=/var/www/biomon/venv/bin/flask sdm worker

# Restart on failure (but not on SIGTERM — graceful stop)
Restart=on-failure
RestartSec=10

# Stop: SIGTERM → the worker finishes the current task, then exits.
# If it does not exit within 5 minutes — SIGKILL.
TimeoutStopSec=300
KillSignal=SIGTERM

# Env file with DATABASE_URI and other secrets
EnvironmentFile=/var/www/biomon/.env

# Logging via journald (unbuffered)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sdm-worker

[Install]
WantedBy=multi-user.target
